# src/config/inittab
::sysinit:/etc/init.d/rcS
::sysinit:/etc/init.d/local.sh
::shutdown:/bin/umount -a -r
::restart:/sbin/init
::respawn:/sbin/getty -L ttymxc0 115200 vt100
